# src/hostproxy/update-hostproxy.cron
# Update hostproxy certs every 8 hours
3 0-23/8 * * * root /opt/hostproxy/host_dist 2>&1 /dev/null

